# tinydonation-supervisor.conf
[program:tinydonation-com]
command = gunicorn runme:app -b localhost:8013
directory = /var/www/tinydonation-com/
user = wwwuser
